# aloha/puppet/zulip/files/cron.d/calculate-first-visible-message-id
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
USER=aloha

0 7 * * * aloha /home/aloha/deployments/current/manage.py calculate_first_visible_message_id --lookback-hours 30
