# aloha/puppet/zulip/files/cron.d/email-mirror
MAILTO=root

* * * * *   aloha cd /home/aloha/deployments/current && ./manage.py email_mirror
