# aloha/puppet/zulip/files/cron.d/invoice-plans
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
USER=aloha

0 22 * * * aloha /home/aloha/deployments/current/manage.py invoice_plans
