# aloha/puppet/zulip_ops/files/cron.d/fetch-contributor-data
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
USER=aloha

0 8 * * * aloha /home/aloha/deployments/current/tools/fetch-contributor-data
