# aloha/puppet/zulip_ops/files/cron.d/test_zephyr_personal_mirrors
# We don't actually need a valid Kerberos cache since these are sent
# unauth anyway -- but a cache is required for zwrite to run.

* * * * * aloha env KRB5CCNAME=/home/aloha/ccache/zmirror-tabbott zwrite -c aloha-mirror-nagios -i nagios-test -m test -Saloha-nagios@mit.edu -d -q >/dev/null 2>/dev/null
